# rc.local
#!/bin/sh
echo "Performing first run tasks..."
/etc/my_init.d/firstrun.sh 

echo "Starting Guacamole..."
/etc/init.d/tomcat7 start
/etc/init.d/guacd start
tail -f /dev/null
